# target/dropwizard-deb-package/files/lib/systemd/system/DiscoveryDiary.service
[Unit]
Description=DiscoveryDiary

[Service]
ExecStart=/usr/bin/DiscoveryDiary
Restart=always
OOMScoreAdjust=-1000
User=dropwizard
Group=dropwizard

[Install]
WantedBy=multi-user.target
